# web/backend/chaoxing-celery.service
# systemd服务配置 - Celery Worker
# 安装位置: /etc/systemd/system/chaoxing-celery.service
# 安装命令: sudo cp chaoxing-celery.service /etc/systemd/system/
# 启动命令: sudo systemctl start chaoxing-celery
# 开机自启: sudo systemctl enable chaoxing-celery
# 查看状态: sudo systemctl status chaoxing-celery
# 查看日志: sudo journalctl -u chaoxing-celery -f

[Unit]
Description=Chaoxing Celery Worker Service
After=network.target chaoxing-backend.service
# 如果使用Redis，添加依赖
# After=redis.service
# Requires=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/chaoxing/web/backend
Environment="PATH=/opt/chaoxing/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/chaoxing"
Environment="C_FORCE_ROOT=true"

# Celery Worker启动命令
ExecStart=/opt/chaoxing/.venv/bin/celery -A celery_app worker \
    --loglevel=info \
    --concurrency=4 \
    --max-tasks-per-child=1000 \
    --time-limit=3600 \
    --soft-time-limit=3300 \
    --logfile=/opt/chaoxing/web/backend/logs/celery.log

# 优雅停止
ExecStop=/bin/kill -s TERM $MAINPID
TimeoutStopSec=60

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=0

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 标准输出
StandardOutput=append:/opt/chaoxing/web/backend/logs/celery_stdout.log
StandardError=append:/opt/chaoxing/web/backend/logs/celery_stderr.log

[Install]
WantedBy=multi-user.target

