# web/backend/chaoxing-backend.service
# systemd服务配置 - FastAPI后端
# 安装位置: /etc/systemd/system/chaoxing-backend.service
# 安装命令: sudo cp chaoxing-backend.service /etc/systemd/system/
# 启动命令: sudo systemctl start chaoxing-backend
# 开机自启: sudo systemctl enable chaoxing-backend
# 查看状态: sudo systemctl status chaoxing-backend
# 查看日志: sudo journalctl -u chaoxing-backend -f

[Unit]
Description=Chaoxing Backend Service (FastAPI)
After=network.target
# 如果使用PostgreSQL，添加依赖
# After=postgresql.service
# Requires=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/chaoxing/web/backend
Environment="PATH=/opt/chaoxing/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/chaoxing"

# 使用gunicorn生产部署（推荐）
ExecStart=/opt/chaoxing/.venv/bin/gunicorn app:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 120 \
    --access-logfile /opt/chaoxing/web/backend/logs/access.log \
    --error-logfile /opt/chaoxing/web/backend/logs/error.log \
    --log-level info

# 或使用uvicorn直接运行（开发）
# ExecStart=/opt/chaoxing/.venv/bin/uvicorn app:app \
#     --host 0.0.0.0 \
#     --port 8000 \
#     --workers 4 \
#     --log-level info

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=0

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全选项（可选，根据需要启用）
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/opt/chaoxing/web/backend/logs /opt/chaoxing/web/backend/data

# 标准输出和错误输出
StandardOutput=append:/opt/chaoxing/web/backend/logs/backend.log
StandardError=append:/opt/chaoxing/web/backend/logs/backend_error.log

[Install]
WantedBy=multi-user.target

